# tomee/tomee-deb/src/main/resources/control/postrm.sh
#!/bin/sh -e

update-rc.d -f tomee-${classifier} remove
update-alternatives --remove tomee /etc/init.d/tomee-${classifier}