# tomee/tomee-deb/src/main/resources/control/postrm.sh
#!/bin/sh

update-rc.d -f tomee-${classifier} remove